# deploy/systemd/monitor-agent.service
[Unit]
Description=Monitor agent
After=network-online.target monitor-server.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/monitor
ExecStart=/opt/monitor/.venv/bin/python -m agent.agent -c /opt/monitor/agent/config.toml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

