# deploy/systemd/monitor-server.service
[Unit]
Description=Monitor server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/monitor
EnvironmentFile=/opt/monitor/server.env
ExecStart=/opt/monitor/.venv/bin/uvicorn server.main:app --host 127.0.0.1 --port 8800
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

